# scripts/whisper-dictation.service
[Unit]
Description=Whisper Dictation Daemon
After=network.target sound.target

[Service]
Type=simple
User=%u
WorkingDirectory=%h/whisper-dictation
ExecStart=%h/whisper-dictation/venv/bin/python -m whisper_dictation.main --daemon
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=GEMINI_API_KEY=${GEMINI_API_KEY}

[Install]
WantedBy=default.target
